# templates/chezmoi-update.timer
# {{ ansible_managed }}

[Unit]
Description="Update chezmoi timer"

[Timer]
Unit=chezmoi-update.service
OnCalendar={{ chezmoi_update_time }}
{% if chezmoi_update_randomized_delay != "" %}
RandomizedDelaySec={{ chezmoi_update_randomized_delay }}
{% endif %}
Persistent=true

[Install]
WantedBy=timers.target
